# scripts/com.pantheon.git-fetch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.pantheon.git-fetch.plist
  每 60 秒執行一次 git fetch，偵測到新 commit 才寫 flag 觸發 upload agent。
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.pantheon.git-fetch</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>check_new_commits.sh</string>
  </array>

  <key>StartInterval</key>
  <integer>60</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>WorkingDirectory</key>
  <string>PLACEHOLDER_HOME/project-pantheon/scripts</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>PLACEHOLDER_HOME</string>
  </dict>
</dict>
</plist>
